# jt-wazuh-mgr.service
[Unit]
Description=JT Wazuh Agent Manager Web UI
Documentation=https://github.com/jasoncheng7115/jt-wazuh-mgr
After=network.target wazuh-manager.service
Wants=wazuh-manager.service

[Service]
Type=simple
WorkingDirectory=/opt/jt-wazuh-mgr
ExecStart=/usr/bin/python3 /opt/jt-wazuh-mgr/wazuh_agent_mgr.py --web --ssl-auto
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jt-wazuh-mgr

[Install]
WantedBy=multi-user.target
